# home/.config/systemd/user/XSession.target
[Unit]
Description=Xorg up and running
BindsTo=graphical-session.target

